# package/base-files/files/etc/init.d/umount
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2015 - 2029. Hydra Project.
#

STOP=90

restart() {
	:
}

stop() {
	sync
	/bin/umount -a -d -r
}
